# factory/schedulers/factory.cron.template
# Factory cron entries — hand-install fallback. Prefer the driver generate
# these lines with the exact paths filled in:
#   node ~/.factory/runtime/factory/driver/factory.mjs schedule --declare --install --project /path/to/project
# If copying by hand into `crontab -e`: fix the project path, and point the
# log redirect at the project's machine-side state dir
# (~/.factory/projects/<name>-<hash8>/ — `ls ~/.factory/projects/` to find it).
# The repo carries no factory logs.
# PATH must include node, claude, and gh (cron gets a minimal PATH).
# The driver is the machine runtime (O6): bootstrap once with
#   git clone <repo-url> ~/.factory/runtime
# --scheduled: doctor preflight — a misconfigured factory aborts loudly
# instead of half-running (NOTES item 25).
PATH=/home/YOU/.local/bin:/usr/local/bin:/usr/bin:/bin
30 8  * * 1-5  node $HOME/.factory/runtime/factory/driver/factory.mjs triage --project /path/to/project --scheduled >> $HOME/.factory/projects/PROJECT-KEY/log/cron.out 2>&1
0  9  * * 1-5  node $HOME/.factory/runtime/factory/driver/factory.mjs dev    --project /path/to/project --scheduled >> $HOME/.factory/projects/PROJECT-KEY/log/cron.out 2>&1
30 13 * * 1-5  node $HOME/.factory/runtime/factory/driver/factory.mjs report --project /path/to/project --scheduled >> $HOME/.factory/projects/PROJECT-KEY/log/cron.out 2>&1
