# factory/schedulers/launchd.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS scheduler for one Factory mode. Make one copy per mode you want:
    com.factory.MYPROJECT.dev.plist, .triage.plist, .report.plist
  Replace __MODE__ (dev|triage|report), __PROJECT__ (absolute path),
  __HOME__ (absolute home dir — launchd does not expand ~),
  __NODE__ (`which node`), __HOUR__/__MINUTE__ (24h start time).
  The driver is the machine runtime (O6): git clone <repo-url> ~/.factory/runtime
  Install:  cp <file> ~/Library/LaunchAgents/ && launchctl load ~/Library/LaunchAgents/<file>
  Trigger now (test):  launchctl start com.factory.MYPROJECT.__MODE__
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.factory.MYPROJECT.__MODE__</string>
    <key>ProgramArguments</key>
    <array>
        <string>__NODE__</string>
        <string>__HOME__/.factory/runtime/factory/driver/factory.mjs</string>
        <string>__MODE__</string>
        <string>--project</string>
        <string>__PROJECT__</string>
        <string>--scheduled</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>__HOUR__</integer>
        <key>Minute</key>
        <integer>__MINUTE__</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>__PROJECT__</string>
    <key>StandardOutPath</key>
    <string>__PROJECT__/.factory/log/launchd-__MODE__.out</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT__/.factory/log/launchd-__MODE__.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <!-- launchd jobs get a minimal PATH; include node's and claude's dirs -->
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
